# systemd/user/phone-status-notify.service
[Unit]
Description=PhoneBridge status notification updater
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/home/jegan/.config/hypr/scripts/phone-status-notify.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
